# srcs/grafana/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command = sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"

[program:grafana]
command = sh -c "grafana-server --config=/etc/grafana.ini --homepath=/usr/share/grafana && kill -s SIGTERM $(cat supervisord.pid)"
